# honey.service
[Unit]
Description=honeypot
Type=simple
Restart=always
ExecStart=/home/odroid/honey.sh

[Install]
WantedBy=multi-user.target