# honey.service
[Unit]
Description=honeypot

[Service]
Type=simple
Restart=always
ExecStart=/usr/bin/python3 $HOME/honey.py

[Install]
WantedBy=multi-user.target
